# supervisor.d/tailscale.ini
[program:tailscale]
command=/usr/local/bin/containerboot
process_name=%(program_name)s
numprocs=1
directory=/tmp
umask=022
priority=999
autostart=true
autorestart=unexpected
startsecs=10
startretries=3
